# templates/postgresql-dump.timer.j2
[Unit]
Description=Postgresql Dump Backup

[Timer]
OnCalendar={{ item.frequency }}

[Install]
WantedBy=multi-user.target

